# Utils/NaoConfigFiles/deploy/v6/etc/systemd/system/brainwasher.service
[Service]
Type=oneshot
ExecStart=/usr/bin/brainwash
