# Utils/NaoConfigFiles/etc/init.d/naoth
#!/sbin/runscript

description="NaoTH init script"

depends(){
    after naoqi
}

start(){
  ebegin "starting naoth"
    naoth start
  eend $?
}

stop(){
  ebegin "stopping naoth"
    naoth stop
  eend $?
}

restart(){
  stop
  sleep 2
  start
}


