# Utils/NaoConfigFiles/deploy/v6/etc/init.d/checkpart-dummy
#!/sbin/runscript
##
## checkpart
##
## Copyright (C) 2010-2012 Aldebaran Robotics
##


depend() {
  provide checkpart
}

start() {
  :
}

stop() {
  :
}

restart() {
  :
}
