# Utils/NaoConfigFiles/deploy/v6/nao/etc/systemd/user/naoth.path
[Unit]
Description=Monitor the file for changes

[Path]
PathExists=/home/nao/Config/nao.info
Unit=naoth.service

[Install]
WantedBy=multi-user.target